# dist/gnutheme.cron.example
#!/bin/sh
#
# crontab script for gnutheme
#

# Change this to the directory where you installed it
prefix=$HOME/gnutheme/
pidfile=var/gnutheme.pid

cd "$prefix" || {
	echo "$0: couldn't cd to $prefix"
	exit 1
}

if [ -f "$pidfile" ] ; then
  PID=`cat "$pidfile"`

  if kill -0 $PID >/dev/null 2>&1; then
    exit 0
  fi

  rm -f "$pidfile"
fi

echo "$0: couldn't find gnutheme running, restarting..."

./bin/gnutheme-services

exit 0
